# templates/dev.mysql.plist.erb
<?xml version="1.0" encoding="UTF-8"?>

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.mysql</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>A pretty okay relational database.</string>

      <key>Port</key>
      <integer><%= @port %></integer>

      <key>Socket</key>
      <string><%= @socket %></string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string><%= @executable %></string>
      <string>--defaults-file=<%= "#{@configdir}/my.cnf" %></string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>ExitTimeOut</key>
    <integer>1</integer>

    <key>UserName</key>
    <string><%= @user %></string>

    <key>WorkingDirectory</key>
    <string><%= @datadir %></string>

    <key>StandardErrorPath</key>
    <string><%= @logdir %>/console.log</string>

    <key>StandardOutPath</key>
    <string><%= @logdir %>/console.log</string>
  </dict>
</plist>
